# deploy/quadlet/beehive-digest.container
[Unit]
Description=Beehive daily digest email

[Container]
Image=localhost/beehive:latest
ContainerName=beehive-digest
Volume=beehive-data.volume:/data
Environment=DB_PATH=/data/beehive.db
Environment=DIGEST_EMAIL_TO=you@example.com
Environment=DIGEST_EMAIL_FROM=beehive@example.com
# ACS connection string for outbound digest/alert email. Create it as a rootless Podman secret
# on the host (see deploy/README.md); it is never baked into the image or committed.
Secret=beehive-acs-connection,type=env,target=ACS_CONNECTION_STRING
Exec=-m scripts.run_collector --mode digest

[Service]
Type=oneshot
CPUQuota=100%
CPUWeight=20
MemoryHigh=256M
MemoryMax=384M
TimeoutStartSec=120
