# deploy/quadlet/beehive-digest.timer
# deploy/quadlet/beehive-digest.timer
[Unit]
Description=Evaluate Beehive email groups hourly (host tz); each group's own send interval decides whether it actually sends

[Timer]
OnCalendar=hourly
Persistent=true
Unit=beehive-digest.service

[Install]
WantedBy=timers.target
