# deploy/quadlet/beehive-digest.timer
# deploy/quadlet/beehive-digest.timer
[Unit]
Description=Evaluate Beehive email group schedules every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=beehive-digest.service

[Install]
WantedBy=timers.target
